# files/etc/service/librenms/run
#!/bin/bash -eu


cd /opt/librenms
git pull
exec php /opt/librenms/poller.php -l -h all