# files/etc/service/librenms/run
#!/bin/bash -eu

cd /opt/librenms
exec php /opt/librenms/poller.php -l -h all